# supervisord.conf
[supervisord]
nodaemon=true

[program:tutum-agent]
command=/usr/bin/tutum-agent

[program:sshd]
command=/usr/sbin/sshd -D

