# tutum-agent/rc.local
#!/bin/bash

TUTUM_TOKEN=${1:-}

set -e

: ${TUTUM_TOKEN:?"-e TUTUM_TOKEN is not set. This is the hex value shown in Bring Your Own Node dialog."}

# make sure cgroup-lite -> cgroup-bin transition takes place
/cgroupfs_mount.sh

curl -Ls https://files.tutum.co/scripts/install-agent.sh | sudo sh -s ${TUTUM_TOKEN} || true

/usr/bin/supervisord

# tutum-agent leaves docker daemon alive
# it's important to kill docker before terminating container
# see https://github.com/jpetazzo/dind/issues/19
/sbin/start-stop-daemon --stop --pidfile "/var/run/docker.pid"

exit 0
